# mtl/crontab
* * * * * /var/www/dida/mtl/email.sh 1 >> /tmp/email.log
